# nctu-hw/roles/external/criecm.common/files/upsmon_openbsd_rc
#!/bin/sh

daemon="/usr/local/sbin/upsmon"

. /etc/rc.d/rc.subr

rc_cmd $1
